# demos/microshift-disconnected/files/usb-automation/99-usb_monitor.rules
ACTION=="add", KERNEL=="sd*[!0-9]", SUBSYSTEM=="block", RUN+="/usr/local/bin/usb_check.sh %E{DEVNAME}"